# .docker/php/cron/crontab
* * * * * cd /var/www && wp cron event run --due-now >> /var/log/cron/cron.log 2>&1
